# etc/cron/crontab
30 2 * * 1 	/opt/letsencrypt/letsencrypt-auto renew >> /var/log/le-renew.log
35 2 * * 1 	nginx -s reload
